# skills/n8n-self-host/references/n8n-update.timer
[Unit]
Description=Run n8n auto-update daily

[Timer]
OnCalendar=*-*-* 04:17:00
Persistent=true
RandomizedDelaySec=600
Unit=n8n-update.service

[Install]
WantedBy=timers.target
